# seafile_server/supervisord.conf
[supervisord]
nodaemon=true

[program:sshd]
command=/usr/sbin/sshd -D

[program:nginx]
command = /usr/sbin/nginx
autostart = true

[program:seafile]
command = /root/seafile_run.sh
